# src/etc/g7ctrl.service.in
[Unit]
Description=GM7 Service Daemon

[Service]
Type=simple
ExecStart=@prefix@/bin/@PACKAGE_NAME@

[Install]
WantedBy=multi-user.target
