# src/etc/g7ctrl.service.in
[Unit]
Description=GM7 Service Daemon

[Service]
ExecStart=@prefix@/bin/@PACKAGE_NAME@ -dn

[Install]
WantedBy=multi-user.target
