# dpkg/debian/preinst
#!/bin/sh
# preinst script for g7ctrl
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <new-preinst> `install'
#        * <new-preinst> `install' <old-version>
#        * <new-preinst> `upgrade' <old-version>
#        * <old-preinst> `abort-upgrade' <new-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    install)        
	# add g7ctrl system user (requires adduser >= 3.34)
        # don't muck around with this unless you KNOW what you're doing
        echo "Creating/updating g7ctrl user account..."
        adduser --system --no-create-home --home /nonexistent --group \
                --gecos "G7CTRL system user" --shell /bin/false \
                --quiet --disabled-password g7ctrl || {
          # adduser failed. Why?
          if getent passwd g7ctrl >/dev/null ; then
             echo "Non-system user g7ctrl found. I will not overwrite a non-system" >&2
             echo "user.  Remove the user and reinstall g7ctrl." >&2
             exit 1
          fi
          # unknown adduser error, simply exit
          exit 1
        }    

    if [ -x /etc/init.d/g7ctrl ]; then
        service g7ctrl stop
        pkill g7ctrl || true
        if pgrep g7ctrl > /dev/null; then sleep 3; pkill -9 g7ctrl || true; fi
        if pgrep g7ctrl > /dev/null; then
            echo "Old g7ctrl daemon should hava stopped but did not. Please kill manually before installation can proceed"
            exit 1
        fi
        rm -f /var/run/g7ctrl.pid
        update-rc.d -f g7ctrl remove
    fi


    ;;

    upgrade)
    # Stop and remove any old daemon
    if [ -x /etc/init.d/g7ctrl ]; then
        service g7ctrl stop
        pkill g7ctrl || true
        if pgrep g7ctrl > /dev/null; then sleep 3; pkill -9 g7ctrl || true; fi
        if pgrep g7ctrl > /dev/null; then
            echo "Old g7ctrl daemon should hava stopped but did not. Please kill manually before installation can proceed"
            exit 1
        fi
        rm -f /var/run/g7ctrl.pid
        update-rc.d -f g7ctrl remove
    fi
    ;;

    abort-upgrade)
    ;;

    *)
        echo "preinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
